# systemd/hugo-update.timer
[Unit]
Description=HUGO Update Timer

[Timer]
OnCalendar=*:0/10
Unit=hugo-update.service

[Install]
WantedBy=timers.target
