# systemd/update.service
[Unit]
Description=Update Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/update.sh

[Install]
WantedBy=multi-user.target
